# telegram-scraper.service
[Unit]
Description=Telegram Scraper Application
Documentation=https://github.com/AbirHasan2005/TelegramScraper
After=network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=docker
Group=docker
WorkingDirectory=/home/telegram-scraper/TelegramScraper

# Environment
EnvironmentFile=/home/telegram-scraper/TelegramScraper/.env
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"

# Service commands
ExecStart=/usr/bin/docker-compose up
ExecStop=/usr/bin/docker-compose down
ExecReload=/usr/bin/docker-compose restart

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=600
StartLimitBurst=3

# Resource limits
MemoryMax=2G
CPUQuota=200%

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=partial
ProtectHome=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-scraper

[Install]
WantedBy=multi-user.target
